# shitcode.service
[Unit]
Description=ShitCode Web App
After=network.target mysql.service

[Service]
WorkingDirectory=/var/www/shitcode/publish
ExecStart=/usr/bin/dotnet /var/www/shitcode/publish/ErrorLogDashboard.Web.dll
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=shitcode
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
# Ensure you set the DeepSeek API Key here or in /etc/environment
# Environment=DeepSeek__ApiKey=your_api_key_here

[Install]
WantedBy=multi-user.target
